# scripts/imagestool1.service
# systemd unit for Folder2Page / imagestool1.
#
# Install:
#   sudo cp scripts/imagestool1.service /etc/systemd/system/imagestool1.service
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now imagestool1
#
# Logs:
#   journalctl -u imagestool1 -f
#
# Clone for additional tools (imagestool2, imagestool3, ...):
#   1. Copy the project to /root/imagestool2
#   2. Copy this file to /etc/systemd/system/imagestool2.service
#   3. Edit WorkingDirectory + Environment=PORT=5017
#   4. systemctl daemon-reload && systemctl enable --now imagestool2

[Unit]
Description=Folder2Page (imagestool1) — Facebook Auto Poster
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/root/imagestool1
ExecStart=/usr/bin/env bash /root/imagestool1/scripts/start.sh
Restart=on-failure
RestartSec=10
Environment=NODE_ENV=production
Environment=PORT=5016
StandardOutput=journal
StandardError=journal
SyslogIdentifier=imagestool1

[Install]
WantedBy=multi-user.target
